# scripts/deploy/meadow-bootstrap.service
[Unit]
Description=Meadow Box Bootstrap (pairing)
Documentation=https://github.com/duquettelogan/meadow
# Need network AND time sync — the pairing API uses HTTPS, and a wrong
# clock breaks cert validation. Pi has no RTC so this matters every boot.
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target
# meadow.service Requires + After this unit so the DNS server only
# starts once /etc/meadow/box.env has a device_id.
# StartLimit* are [Unit]-level directives in systemd 230+. Previously
# placed in [Service] where systemd silently ignored them.
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=oneshot
RemainAfterExit=yes
User=meadow
Group=meadow

# Working dir + state.
WorkingDirectory=/opt/meadow
Environment=NODE_ENV=production
Environment=BOX_ENV_FILE=/etc/meadow/box.env
Environment=PAIRING_CODE_FILE=/var/lib/meadow/pairing-code
EnvironmentFile=-/etc/meadow/bootstrap.env

# Run via ts-node in dev; for prod, point this at the compiled JS instead.
ExecStart=/usr/bin/env npx ts-node /opt/meadow/src/box/bootstrap.ts

# /etc/meadow must exist and be writable by the meadow user
# (install.sh creates it with mode 0700 owner meadow:meadow).
# /var/lib/meadow holds the box-generated pairing-code file (created
# by install.sh so the unprivileged meadow user can write it).
ReadWritePaths=/etc/meadow /var/lib/meadow

# Bind port 80 for the LAN-facing pairing web page (http://meadow.local)
# without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
# MemoryDenyWriteExecute removed — V8 needs W^X transitions for JIT;
# enabling it sigtraps Node on the first JIT compile (~80ms after
# boot). Confirmed on Pi 5 / Pi OS Lite 64-bit (May 2026 alpha test).

# If pairing fails (code expired, network out), let systemd retry a few
# times before giving up.
Restart=on-failure
RestartSec=15s

[Install]
WantedBy=multi-user.target
