# scripts/deploy/meadow.service
[Unit]
Description=Meadow — kid-safe network filtering
Documentation=https://github.com/duquettelogan/meadow
# time-sync.target: the Pi has no RTC. Booting before NTP has settled
# leaves the clock at epoch, which breaks HTTPS cert validation against
# upstream feeds and the API. Wait for time to come in before starting.
After=network-online.target time-sync.target meadow-bootstrap.service postgresql.service redis-server.service
Wants=network-online.target time-sync.target
# meadow-bootstrap is a Type=oneshot that pairs the box and writes
# /etc/meadow/state.json. We require it so the DNS path only starts
# once a device_id exists (otherwise blocks would have no child to
# attribute to).
Requires=meadow-bootstrap.service postgresql.service redis-server.service
# StartLimit* are [Unit]-level directives in systemd 230+. Previously
# placed in [Service] where they were silently ignored (no rate limit
# applied → flapping restarts under a hard fault).
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=meadow
Group=meadow
WorkingDirectory=/opt/meadow
EnvironmentFile=/etc/meadow/meadow.env
ExecStart=/usr/bin/npx ts-node src/index.ts
Restart=on-failure
RestartSec=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meadow

# Allow binding to port 53 without root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/opt/meadow

# Resource limits
LimitNOFILE=65536
MemoryHigh=512M
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
